# onserver/.config/systemd/user/on_server.service
[Unit]
Description=Onserver service

[Service]
Type=Simple
ExecStart=/usr/bin/python3 /home/onserver/on_server.py
